# dot_config/systemd/user/hyde-wallpaper.service
[Unit]
Description=HyDE Wallpaper Cycler

[Service]
Type=oneshot
ExecStart=hyde-shell wallpaper --random
